# pi/duoclock-monitor.service
# ===========================================================================
# 🧠 DuoClock Monitor — systemd service unit
# ===========================================================================
# The core daemon that reads serial from the ESP32, manages timers,
# controls LEDs, and writes event logs.
#
# Management:
#   systemctl status duoclock-monitor    # Check status
#   journalctl -u duoclock-monitor -f    # View live logs
#   sudo systemctl restart duoclock-monitor  # Restart
# ===========================================================================

[Unit]
Description=DuoClock Button Monitor
After=multi-user.target
# StartLimitIntervalSec=0 allows unlimited restarts (useful for USB reconnects)
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/duoclock
ExecStart=/usr/bin/python3 /opt/duoclock/duoclock_monitor.py
# Restart always — handles USB disconnects, crashes, etc.
Restart=always
RestartSec=5
# Log to journald (viewable via journalctl)
StandardOutput=journal
StandardError=journal

# 🔒 Security hardening
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
