# pi/duoclock-web.service
# ===========================================================================
# 🌐 DuoClock Web Dashboard — systemd service unit
# ===========================================================================
# HTTP server on port 8080 showing 7-day timer activity.
# Reads /var/log/duoclock.log written by duoclock-monitor.
#
# URL: http://clockradio.belairmoon.au:8080/
#
# Management:
#   systemctl status duoclock-web          # Check status
#   journalctl -u duoclock-web -f          # View live logs
#   sudo systemctl restart duoclock-web    # Restart
# ===========================================================================

[Unit]
Description=DuoClock Web Dashboard
# Start after network (needs TCP) and after monitor (needs log file)
After=network.target duoclock-monitor.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/python3 /opt/duoclock/duoclock_web.py
# Restart always — resilient to crashes
Restart=always
RestartSec=5
# Log to journald (viewable via journalctl)
StandardOutput=journal
StandardError=journal

# 🔒 Security hardening
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
